# debian/postinst
#!/bin/sh
# postinst script for zenloadbalancer
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
	#upgrading global.conf
	/usr/local/zenloadbalancer/app/checkglobalconf/checkglobalconf.cgi
	# stop gdnsd
	update-rc.d gdnsd disable > /dev/null
	service gdnsd stop
	# hack snmpd
	if [ ! -f /etc/snmp/snmpd.conf.zensave ]; then
		mv /etc/snmp/snmpd.conf /etc/snmp/snmpd.conf.zensave
		echo "agentAddress udp:127.0.0.1:161" > /etc/snmp/snmpd.conf
		echo "rocommunity public 127.0.0.1/8" >> /etc/snmp/snmpd.conf
		echo "# zen modified: see old config in snmpd.conf.zensave" >> /etc/snmp/snmpd.conf
	fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
